# dot_config/systemd/user/freshclam.service
[Unit]
Description=Updates ClamAV virus definitions

[Service]
Type=oneshot
ExecStart=/usr/bin/freshclam
