# dot_config/systemd/user/home-manager-optimization.timer
[Unit]
Description=Daily nix store optimization

[Timer]
OnCalendar=*-*-* 5:00:00
Persistent=true

[Install]
WantedBy=timers.target
