# dot_config/systemd/user/home-manager-optimization.service
[Unit]
Description=Optimize nix store

[Service]
Type=oneshot
ExecStart=%h/.nix-profile/bin/nix-store --optimise
