# dot_config/systemd/user/maldet-update.service
[Unit]
Description=Updates maldet signatures

[Service]
Type=oneshot
ExecStart=/usr/bin/maldet --update-sigs --force
