# dot_config/systemd/user/pkgfile-update.service
[Unit]
Description=Updates pkgfile database

[Service]
Type=oneshot
ExecStart=/usr/bin/pkgfile --update
