# dot_config/systemd/user/tldr-update.service
[Unit]
Description=Updates tldr pages cache

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c '$$(/usr/bin/which tldr) --clear-cache --update'
