# dot_config/systemd/user/tldr-update.service
[Unit]
Description=Updates tldr pages cache

[Service]
Type=oneshot
ExecStart=/usr/bin/tldr --clear-cache --update
