# dot_config/systemd/user/trash-empty.service
[Unit]
Description=Deletes trash items older than 30 days

[Service]
Type=oneshot
ExecStart=%h/.nix-profile/bin/trash-empty 30
