# dot_config/systemd/user/trash-empty.service.tmpl
[Unit]
Description=Deletes trash items older than 30 days

[Service]
Type=oneshot
ExecStart={{ .chezmoi.homeDir }}/{{ .paths.nix_user_profile }}/bin/trash-empty 30
