# dot_config/systemd/user/trash-empty.service.tmpl
[Unit]
Description=Deletes trash items older than 30 days

[Service]
Type=oneshot
ExecStart={{ env "NIX_USER_PROFILE" }}/bin/trash-empty 30
