# dot_config/exact_systemd/user/trash-empty.timer
[Unit]
Description=Delete old trash items

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
